# scripts/askme-dds-bridge.service
[Unit]
Description=askme DDS bridge (Unix socket relay for robot telemetry)
After=network.target
Wants=network.target

[Service]
Type=simple
User=sunrise
WorkingDirectory=/home/sunrise/askme
EnvironmentFile=-/home/sunrise/askme/.env
ExecStart=/bin/bash -c 'source /opt/ros/humble/setup.bash && python3 /home/sunrise/askme/scripts/poc_dds_bridge.py'
Restart=on-failure
RestartSec=3
StandardOutput=append:/home/sunrise/askme/logs/dds-bridge.log
StandardError=append:/home/sunrise/askme/logs/dds-bridge.log

[Install]
WantedBy=multi-user.target
